# assets/supervisor/caddy.ini
[program:caddy]
command=/usr/bin/start-caddy.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
